# homebrew.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>Set GitHub API token for Homebrew</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/launchctl</string>
		<string>setenv</string>
		<string>HOMEBREW_GITHUB_API_TOKEN</string>
		<string>${HOMEBREW_GITHUB_API_TOKEN}</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
