# systemd/keepass-sync-to-dropbox.path
[Unit]
Description=Sync KeePass changes to Dropbox
Wants=network-online.target
Before=keepass-sync-from-dropbox.service

[Path]
PathChanged=%h/Asiakirjat/KeePass/
MakeDirectory=true
Unit=keepass-sync-to-dropbox.service

[Install]
WantedBy=paths.target

